# wifi-portal/scripts/wifi-check.service
[Unit]
Description=WiFi Check and Hotspot Fallback
After=network.target

[Service]
ExecStart=/usr/local/bin/wifi-check.sh
Restart=always

[Install]
WantedBy=multi-user.target
